# tv.jibb.weatherdesktop.cacheflush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>tv.jibb.weatherdesktop.cacheflush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/blake/Developer/weatherdesktop/flush_wallpaper_cache.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/blake/Library/Logs/weatherdesktop_cacheflush.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/blake/Library/Logs/weatherdesktop_cacheflush.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
<!--
# Install as user LaunchAgent (no sudo needed)
cp /Users/blake/Developer/weatherdesktop/tv.jibb.weatherdesktop.cacheflush.plist ~/Library/LaunchAgents/

# Load into user session
launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/tv.jibb.weatherdesktop.cacheflush.plist

# Check status
launchctl print gui/$(id -u)/tv.jibb.weatherdesktop.cacheflush

# Test run immediately
launchctl kickstart gui/$(id -u)/tv.jibb.weatherdesktop.cacheflush

# Unload
launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/tv.jibb.weatherdesktop.cacheflush

# Remove old system daemon if previously installed
sudo launchctl bootout system/tv.jibb.weatherdesktop.cacheflush
sudo rm /Library/LaunchDaemons/tv.jibb.weatherdesktop.cacheflush.plist
-->
